# flasher/overlays/installer/etc/systemd/system/installer.service
[Unit]
Description=Installer

[Service]
Type=idle
ExecStart=/usr/local/bin/installer.sh
StandardInput=tty
StandardOutput=tty
TTYPath={INSTALLER_TTY}
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=installer.target
