# flasher/overlays/installer/etc/systemd/system/installer.target
[Unit]
Description=Installer target
Documentation=man:systemd.special(7)
Requires=multi-user.target
Conflicts=rescue.service rescue.target
After=basic.target rescue.service rescue.target
AllowIsolate=yes
