# overlays/getty-autologin/etc/systemd/system/serial-getty@ttyS1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin demo --noclear %I $TERM
